# oati-ts-monitor-daemon.service
[Unit]
Description=oati-ts-monitor daemon (Modbus polling + HTTP + logging)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/PATH/TO/oati-ts-monitor
ExecStart=/usr/bin/python3 /PATH/TO/oati-ts-monitor/daemon.py \
  --config /PATH/TO/oati-ts-monitor/monitor_config.json \
  --poll 1 \
  --ipc-bind 127.0.0.1 --ipc-port 9102 \
  --http --http-bind 0.0.0.0 --http-port 8080 \
  --log-file /var/log/oati-ts-monitor/monitor.jsonl --log-interval 1

Restart=always
RestartSec=2

# Helpful for debugging: view logs with `journalctl -u oati-ts-monitor-daemon.service -f`
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
