# debian/oati-ts-monitor/lib/systemd/system/oati-ts-monitor.service
[Unit]
Description=oati-ts-monitor daemon (Modbus polling + HTTP + logging)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oati-ts-monitor
Group=oati-ts-monitor
WorkingDirectory=/var/lib/oati-ts-monitor
ExecStart=/usr/bin/oati-ts-monitor-daemon \
  --config /etc/oati-ts-monitor/monitor_config.json \
  --poll 1 \
  --ipc-bind 127.0.0.1 --ipc-port 9102 \
  --http --http-bind 0.0.0.0 --http-port 8080 \
  --log-file /var/log/oati-ts-monitor/monitor.jsonl --log-interval 1

Restart=always
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
